# scripts/com.cluely.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cluely.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/romirpatel/apx/.venv/bin/python</string>
        <string>/Users/romirpatel/apx/backend/main.py</string>
        <string>--daemon</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/romirpatel/apx/backend</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/romirpatel/apx/backend/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/romirpatel/apx/backend/daemon_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/romirpatel/apx/.venv/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/romirpatel/apx/backend</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
